# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql14/bin/postgres -D /data/postgres
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/postgres.log
stderr_logfile=/data/logs/postgres.log

[program:nextjs]
command=node_modules/.bin/next start apps/site -p 3000
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/data/logs/nextjs.log
stderr_logfile=/data/logs/nextjs.log
environment=NODE_ENV="production",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.log
